# conf/etc/service/carbon/run
#!/bin/sh

exec sudo -u telemetry python3 /opt/graphite/bin/carbon-cache.py start --debug 2>&1
