# conf/etc/service/carbon-aggregator/run
#!/bin/sh

exec sudo -u telemetry python3 /opt/graphite/bin/carbon-aggregator.py start --debug 2>&1
